# Server backup/boot/config/plugins/dynamix/mover.cron
# Generated mover schedule:
40 3 * * * /usr/local/sbin/mover |& logger

