# Server backup/boot/config/plugins/dynamix/parity-check.cron
# Generated parity check schedule:
0 1 1 * * /usr/local/sbin/mdcmd check  &> /dev/null

