# Server backup/boot/config/plugins/dynamix/status-check.cron
# Generated array status check schedule:
20 0 * * * /usr/local/emhttp/plugins/dynamix/scripts/statuscheck &> /dev/null

